# 20151230_D-link_DWR-118_V1.01b01/_20151230_D-link_DWR-118_V1.01b01.bin.extracted/squashfs-root-0/usr/etc/inittab
# < id >:< runlevels >:< action >:< process >
::sysinit:/bin/mount -a
ttyS1::sysinit:/usr/etc/rcS
ttyS1::respawn:-/bin/ash

